# launchagents/com.drobdi.daily-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.drobdi.daily-sync</string>

    <!-- __DROBDI_ROOT__ is substituted at install time by scripts/install-launchagents.sh -->
    <!-- /bin/bash -l loads the user profile (PATH includes /usr/local/bin for docker).    -->
    <!-- If launchd fails with "Operation not permitted", grant Full Disk Access to         -->
    <!-- /bin/bash in System Settings → Privacy & Security → Full Disk Access.             -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>__DROBDI_ROOT__/scripts/daily-sync-all.sh</string>
    </array>

    <!-- Daily at 00:00 local time. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/drobdi-daily-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/drobdi-daily-sync.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
