# launchagents/com.drobdi.memory-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.drobdi.memory-sync</string>

    <!-- Syncs drobdi-memory + drobsidian every hour via docker exec (no FDA required).  -->
    <!-- On failure, sends a macOS notification via osascript.                           -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/usr/local/bin/docker exec openclaw-gateway bash /home/node/.openclaw/scripts/daily-sync-container.sh || osascript -e 'display notification "Vérifie /tmp/drobdi-memory-sync.log" with title "drobdi sync failed"'</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/tmp/drobdi-memory-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/drobdi-memory-sync.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
