# launchagents/com.drobdi.memory-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.drobdi.memory-sync</string>

    <!-- Syncs drobdi-memory + drobsidian every hour via docker exec (no FDA required).  -->
    <!-- Skips silently if Docker Desktop is not running (e.g. after reboot).            -->
    <!-- On failure, sends a macOS notification via osascript.                           -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>echo "[$(date '+%Y-%m-%d %H:%M:%S')] sync started"; /usr/local/bin/docker info &gt;&amp;/dev/null || { echo "[$(date '+%Y-%m-%d %H:%M:%S')] Docker not running — skipping"; exit 0; }; /usr/local/bin/docker exec openclaw-gateway bash /home/node/.openclaw/scripts/daily-sync-container.sh; code=$?; echo "[$(date '+%Y-%m-%d %H:%M:%S')] sync finished (exit $code)"; [ $code -ne 0 ] &amp;&amp; osascript -e 'display notification "Vérifie /tmp/drobdi-memory-sync.log" with title "drobdi sync failed"'; exit $code</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/tmp/drobdi-memory-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/drobdi-memory-sync.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
